# sidecar/mautrix-meta-shim/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:synapse]
command=python -m synapse.app.homeserver --config-path /data/homeserver.yaml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:mautrix-meta]
command=/usr/local/bin/mautrix-meta -c /data/mautrix-meta/config.yaml
directory=/data/mautrix-meta
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
; Wait 10s for Synapse to come up first
startsecs=10

[program:shim]
command=/usr/local/bin/shim
environment=FBM_PORT="%(ENV_FBM_PORT)s",FBM_AUTH_TOKEN="%(ENV_FBM_AUTH_TOKEN)s",FBM_HMAC_SECRET="%(ENV_FBM_HMAC_SECRET)s",FBM_WEBHOOK_URL="%(ENV_FBM_WEBHOOK_URL)s",SYNAPSE_ADMIN_TOKEN="%(ENV_SYNAPSE_ADMIN_TOKEN)s",SYNAPSE_URL="http://localhost:8008",BRIDGE_BOT_MXID="@metabot:fbm.local"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
startsecs=15
